# tools/vault-gc.service.in
[Unit]
Description=Vault garbage collection

[Service]
ExecStart=@prefix@/libexec/vault/gc-default
